# supervisord.conf
[supervisord]
nodaemon true

[program:ipfs]
command=/go/bin/ipfs daemon